# docker/s6-overlay/s6-rc.d/queue-worker/run
#!/command/execlineb -P
# Laravel Queue Worker - Processes queue jobs
# Controlled by QUEUE_WORKER_ENABLED environment variable
cd /var/www/html
/usr/local/bin/php artisan queue:work --tries=3 --timeout=90
