# docker/s6-overlay/s6-rc.d/queue-worker/run
#!/command/execlineb -P
# Laravel Queue Worker - Processes background jobs
# Always active in both dev and production
s6-setuidgid www-data
cd /var/www/html
/usr/local/bin/php artisan queue:work --tries=3 --timeout=90
