# docker/s6-overlay/s6-rc.d/init-assets/up
s6-setuidgid www-data /etc/s6-overlay/scripts/init-assets.sh
