# docker/s6-overlay/s6-rc.d/init-usermod/up
/etc/s6-overlay/scripts/init-usermod.sh
