# docker/s6-overlay/s6-rc.d/scheduler/run
#!/command/execlineb -P
# Laravel Scheduler - Runs scheduled tasks
# Only active in production mode (disabled in local/development)
s6-setuidgid www-data
cd /var/www/html
/usr/local/bin/php artisan schedule:work
