# docker/s6-overlay/s6-rc.d/vite-dev/run
#!/command/with-contenv sh
# Vite Dev Server - Hot Module Replacement
# DEVELOPMENT: runs npm run dev
# PRODUCTION: exits immediately (assets built by init-assets)

if [ "$APP_ENV" = "production" ]; then
    echo "vite-dev: disabled in production"
    exec s6-svc -Od .
fi

echo "vite-dev: starting npm run dev"
cd /var/www/html
exec s6-setuidgid www-data npm run dev
