# docker/s6-overlay/s6-rc.d/vite-dev/run
#!/command/execlineb -P
# Vite Dev Server - Hot Module Replacement
# Only runs in development (removed from bundle in production by init-usermod)

s6-setuidgid www-data
cd /var/www/html
/usr/bin/npm run dev
