# docker/s6-overlay/s6-rc.d/vite-dev/run
#!/command/execlineb -P
# Vite Dev Server - Hot Module Replacement for development
# Only active in development mode (APP_ENV=local)
cd /var/www/html
/usr/bin/npm run dev
