# docker/s6-overlay/s6-rc.d/nutgram/run
#!/command/execlineb -P
s6-setuidgid www-data
cd /var/www/html
/usr/local/bin/php artisan nutgram:listen
