# docker/s6-overlay/s6-rc.d/nginx/run
#!/command/execlineb -P
# Nginx - Web Server
/usr/sbin/nginx -g "daemon off;"
