# docker/s6-overlay/s6-rc.d/php-fpm/run
#!/command/execlineb -P
# PHP-FPM - FastCGI Process Manager
/usr/local/sbin/php-fpm -F
